# deploy/cloud/systemd/ibgateway-restart.timer
[Unit]
Description=Restart IB Gateway nightly at 23:45 ET

[Timer]
OnCalendar=*-*-* 23:45 America/New_York
Persistent=false

[Install]
WantedBy=timers.target
